# cosmic-wayland/udev/99-magic-mouse-gestures.rules
# Magic Mouse Linux Gestures
# Gives the active graphical user access to Apple Magic Mouse hidraw nodes.
# Reconnect the mouse after installing this rule.
#
# Product IDs:
#   030d = Magic Mouse
#   0269 = Magic Mouse 2 Lightning
#   0323 = Magic Mouse 2 USB-C / newer Magic Mouse variants
#
# The KERNELS matches are important for Bluetooth hidraw nodes such as:
#   /sys/.../0005:05AC:030D.000E/hidraw/hidraw8
# ENV{HID_ID} is not always available on the hidraw child itself.

ACTION!="add|change", GOTO="magic_mouse_gestures_end"
SUBSYSTEM!="hidraw", GOTO="magic_mouse_gestures_end"
KERNEL!="hidraw*", GOTO="magic_mouse_gestures_end"

# Mark Apple Magic Mouse hidraw nodes.
KERNELS=="0005:05AC:030D.*", ENV{MAGIC_MOUSE_GESTURES}="1"
KERNELS=="0005:05AC:0269.*", ENV{MAGIC_MOUSE_GESTURES}="1"
KERNELS=="0005:05AC:0323.*", ENV{MAGIC_MOUSE_GESTURES}="1"
KERNELS=="0003:05AC:030D.*", ENV{MAGIC_MOUSE_GESTURES}="1"
KERNELS=="0003:05AC:0269.*", ENV{MAGIC_MOUSE_GESTURES}="1"
KERNELS=="0003:05AC:0323.*", ENV{MAGIC_MOUSE_GESTURES}="1"

# Some stacks expose HID_ID directly.
ENV{HID_ID}=="0005:000005AC:0000030D", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_ID}=="0005:000005AC:00000269", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_ID}=="0005:000005AC:00000323", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_ID}=="0003:000005AC:0000030D", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_ID}=="0003:000005AC:00000269", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_ID}=="0003:000005AC:00000323", ENV{MAGIC_MOUSE_GESTURES}="1"

# USB parent fallback.
ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="030d", ENV{MAGIC_MOUSE_GESTURES}="1"
ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="0269", ENV{MAGIC_MOUSE_GESTURES}="1"
ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="0323", ENV{MAGIC_MOUSE_GESTURES}="1"

# Product-name fallback, including intentionally mixed-case Bluetooth names.
ATTRS{name}=="*Magic Mouse*", ENV{MAGIC_MOUSE_GESTURES}="1"
ATTRS{name}=="*MaGiK MoUsE*", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_NAME}=="*Magic Mouse*", ENV{MAGIC_MOUSE_GESTURES}="1"
ENV{HID_NAME}=="*MaGiK MoUsE*", ENV{MAGIC_MOUSE_GESTURES}="1"

# Generic access: logind active-session ACL + magicmouse group fallback.
ENV{MAGIC_MOUSE_GESTURES}=="1", TAG+="uaccess", GROUP="magicmouse", MODE="0660"

# Optional installer-generated per-user ACL fallback. install.sh replaces
# __MAGIC_MOUSE_USER__ with the installing user; if this placeholder remains,
# the RUN line is harmless and intentionally disabled.
ENV{MAGIC_MOUSE_GESTURES}=="1", ENV{MAGIC_MOUSE_ACL_USER}="__MAGIC_MOUSE_USER__"
ENV{MAGIC_MOUSE_GESTURES}=="1", ENV{MAGIC_MOUSE_ACL_USER}!="__MAGIC_MOUSE_USER__", RUN+="/usr/bin/setfacl -m u:%E{MAGIC_MOUSE_ACL_USER}:rw %E{DEVNAME}"

LABEL="magic_mouse_gestures_end"
